# docs/marque.service
[Unit]
Description=Marque Service
After=network.target

[Service]
Type=simple
ExecStart=/home/user/bin/marque -c /home/user/.config/marque/config.json serve
WorkingDirectory=/home/user
User=user
Group=user
StandardOutput=syslog
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
